# supervisord.conf
[include]
files = /home/agent/supervisor.d/*.conf

[unix_http_server]
file=/tmp/supervisor.sock  ; No auth needed — socket is container-internal only

[supervisord]
nodaemon=true
logfile=/atlas/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/atlas/logs

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:init]
command=/atlas/app/init.sh
autorestart=false
startsecs=0
startretries=0
priority=10
stdout_logfile=/atlas/logs/init.log
stderr_logfile=/atlas/logs/init-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/atlas/logs/nginx.log
stderr_logfile=/atlas/logs/nginx-error.log

[program:web-ui]
command=bun run /atlas/app/web-ui/index.ts
autostart=false
autorestart=true
priority=40
stdout_logfile=/atlas/logs/web-ui.log
stderr_logfile=/atlas/logs/web-ui-error.log

[program:metrics]
command=bun run /atlas/app/metrics/server.ts
autostart=true
autorestart=true
priority=15
stdout_logfile=/atlas/logs/metrics.log
stderr_logfile=/atlas/logs/metrics-error.log

[program:supercronic]
command=supercronic /home/agent/crontab
autostart=false
autorestart=true
priority=60
stdout_logfile=/atlas/logs/cron.log
stderr_logfile=/atlas/logs/cron-error.log
